# system/ai-login-check.service
[Unit]
Description=Audio Interface Login Check for already connected devices
After=graphical-session.target
# Additional delay after graphical login
Wants=graphical-session.target

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 10
ExecStart=/usr/local/bin/ai-login-check.sh
RemainAfterExit=no
Environment=DISPLAY=:1

[Install]
WantedBy=default.target
